# deploy/systemd/lark-ai-task-queue.service
# systemd service 模板(Linux 服务器,opt-in)。
# 放到 ~/.config/systemd/user/ (用户级,推荐) 或 /etc/systemd/system/ (系统级)。
# 替换占位符:
#   __ROOT__      项目绝对路径
#   __NODE_BIN__  node/lark-cli/claude 所在目录(`dirname "$(command -v node)"`)
# 配合同目录的 .timer 使用。
[Unit]
Description=Lark AI Task Queue — 跑一轮飞书 AI 任务队列
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__ROOT__
Environment=LARK_CLI_NO_PROXY=1
# 关键:systemd 的 PATH 很干净,必须把 node/lark-cli/claude 目录放进来,否则找不到命令。
Environment=PATH=__NODE_BIN__:/usr/local/bin:/usr/bin:/bin
ExecStart=__NODE_BIN__/node __ROOT__/bin/larkaq run
